# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

; ─────────────────────────────────────────────────────────────────
; PostgreSQL 15
; ─────────────────────────────────────────────────────────────────
[program:postgresql]
command=/usr/lib/postgresql/15/bin/postgres -D /var/lib/postgresql/data/pgdata
user=postgres
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/supervisor/postgres.log
stderr_logfile=/var/log/supervisor/postgres.log
stopsignal=QUIT
stopwaitsecs=60

; ─────────────────────────────────────────────────────────────────
; Redis 7
; ─────────────────────────────────────────────────────────────────
[program:redis]
command=redis-server --save 60 1 --loglevel warning --maxmemory 512mb --maxmemory-policy allkeys-lru --bind 127.0.0.1
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/supervisor/redis.log
stderr_logfile=/var/log/supervisor/redis.log

; ─────────────────────────────────────────────────────────────────
; Backend — FastAPI (waits for pg/redis via entrypoint)
; ─────────────────────────────────────────────────────────────────
[program:backend]
command=uvicorn main:app --host 0.0.0.0 --port 8000
directory=/app/backend
autostart=true
autorestart=true
priority=30
stdout_logfile=/var/log/supervisor/backend.log
stderr_logfile=/var/log/supervisor/backend.log

; ─────────────────────────────────────────────────────────────────
; Scheduler — APScheduler cron jobs
; ─────────────────────────────────────────────────────────────────
[program:scheduler]
command=python scheduler.py
directory=/app/scheduler
autostart=true
autorestart=true
priority=40
stdout_logfile=/var/log/supervisor/scheduler.log
stderr_logfile=/var/log/supervisor/scheduler.log

; ─────────────────────────────────────────────────────────────────
; Frontend — Next.js standalone
; ─────────────────────────────────────────────────────────────────
[program:frontend]
command=node server.js
directory=/app/frontend
autostart=true
autorestart=true
priority=40
environment=NODE_ENV="production",HOSTNAME="0.0.0.0",PORT="12345"
stdout_logfile=/var/log/supervisor/frontend.log
stderr_logfile=/var/log/supervisor/frontend.log
